# install/modules/LDAP/ldap.cron
45       22      *       *       *	root       fazbkp.sh >> /var/log/backup-ok.log 2> /var/log/backup-erro.log 
